# scripts/fticks-upstart.conf
start on startup
start on runlevel [2345]
stop on runlevel [016]

respawn

script
  chdir [FLOG ROOT]
  . env/bin/activate
  exec env/bin/python scripts/fticks.py -d -u [URL] -p [PIPE]
end script
